# tce/mydata/etc/periodic/hourly/logrotate
#!/bin/sh
# -*- mode: sh -*-

/usr/local/sbin/logrotate -v /etc/logrotate.conf 2>&1
