# tce/mydata/etc/periodic/daily/time
#!/bin/sh

NTP_SERVER=us.pool.ntp.org

ntpdate $NTP_SERVER | /usr/bin/logger