# tce/mydata/etc/periodic/daily/time
#!/bin/sh

ntpdate tick.ucla.edu | /usr/bin/logger
