# tce/mandatory/etc/periodic/hourly/logrotate
#!/bin/sh
# -*- mode: sh -*-

/usr/local/sbin/logrotate /etc/logrotate.conf 2>&1
