# tce/mandatory/etc/periodic/hourly/morningstar
#!/bin/sh
# -*- mode: sh -*-

/app/morningstar --driver pl2303 --syslog morningstar > /dev/null
