# tce/mydata/opt/glacier/crontab
*/5 * * * * root /opt/glacier/mirror-obsidian | /usr/bin/logger -t mirror-obsidian
